# etc/s6/x11vnc/run
#!/bin/execlineb -WP
foreground { sleep 2 }
/sbin/su-exec admin
/usr/bin/x11vnc -display :0
