# etc/s6/openbox/run
#!/bin/execlineb -WP
foreground { sleep 1 }
/sbin/su-exec admin
/usr/bin/openbox-session
